# overlay/etc/services.d/nginx-access-log/run
#!/usr/bin/execlineb -P

#======================================================================================================================
# From https://github.com/just-containers/nginx-loadbalancer/blob/master/rootfs/etc/services.d/nginx-access-log/run
#======================================================================================================================

echo " - starting Nginx Access Log"
logutil-service -f /var/run/s6/nginx-access-log-fifo /var/log/nginx-access-logs
