# overlay/etc/s6-overlay/s6-rc.d/nginx-access-log/run
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x nginx-access-log

# Forward nginx access log to Docker
def main [...args] {
    # if forwading is not enable, disable the service
    if not (bf env check NGINX_FORWARD_ACCESS_LOG) {
        bf-s6 svc down nginx-access-log
        return
    }

    # forward the access log file to nginx
    bf-s6 log forward nginx-access-log (bf env NGINX_ACCESS_LOG)
}
