# overlay/etc/periodic/daily/nginx-reload
#!/bin/sh

set -euo pipefail


#======================================================================================================================
# Use executable to reload Nginx.
#======================================================================================================================

nginx-reload
