# overlay/etc/periodic/daily/nginx-reload
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Use executable to reload Nginx every day
#======================================================================================================================

nginx-reload
