# overlay/etc/s6-overlay/s6-rc.d/nginx/finish
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x nginx

# Show helpful log message
def main [...args] {
    # ensure nginx is closed gracefully
    # Flags:
    #   -s  send a signal to the master process (quit: shut down server gracefully)
    bf write debug "Quitting nginx."
    ^nginx -s quit

    # finish and terminate container
    bf-s6 svc finish --terminate
}
