# overlay/etc/services.d/nginx-error-log/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Forward Nginx error log to Docker
#======================================================================================================================

bcg-forward nginx-error-log /var/log/nginx/error.log
