# overlay/etc/services.d/nginx-error-log/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Forward nginx error log to Docker
#======================================================================================================================

if [[ -e /var/log/nginx/error.log ]]; then
    _echo "Forwarding Nginx error log to docker"
    # -F = follow and keep retrying
    exec tail -F /var/log/nginx/error.log
else
    # wait for Nginx service to start and create log file
    exec sleep 5
fi
