# overlay/etc/services.d/nginx-error-log/run
#!/usr/bin/with-contenv sh

set -e


#======================================================================================================================
# Forward nginx error log to Docker
#======================================================================================================================

if [[ -e /var/log/nginx/error.log ]]; then
    _echo "Forwarding Nginx error log to docker"
    exec tail -F /var/log/nginx/error.log
else
    exec sleep 10
fi
