# overlay/etc/s6-overlay/s6-rc.d/nginx-error-log/run
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x nginx-error-log

# Forward nginx error log to Docker
def main [...args] { bf-s6 log forward nginx-error-log (bf env NGINX_ERROR_LOG) }
