# overlay/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Use base executable to terminate all services
#   -s  send a signal to the master process (quit: shut down server gracefully)
#======================================================================================================================

nginx -s quit
_terminate
