# overlay/etc/services.d/nginx/finish
#!/usr/bin/execlineb -P

#======================================================================================================================
# From https://github.com/just-containers/nginx-loadbalancer/blob/master/rootfs/etc/services.d/nginx/finish
#======================================================================================================================

s6-svscanctl -t /var/run/s6/services
