# overlay/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

#======================================================================================================================
# From https://github.com/just-containers/nginx-loadbalancer/blob/master/rootfs/etc/services.d/nginx/run
#======================================================================================================================

nginx -g "daemon off;"
